# roles/base/prometheus/alertmanager/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus server"

start on filesystem and started docker

pre-start script
    # prometheus_alertmanager_container
    if ! docker inspect {{ prometheus_alertmanager_container }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --volume /var/lib/iptables/debian-jessie:/iptables:ro \
            --volume /etc/iptables-docker/{{ prometheus_alertmanager_container }}:/etc/docker-iptables:ro \
            --volume /dev/log:/dev/log \
            --env CONTAINER_NAME={{ prometheus_alertmanager_container }} \
            --cap-add NET_ADMIN \
            --name {{ prometheus_alertmanager_container }} \
            --volume /etc/container/{{ prometheus_alertmanager_container }}/alertmanager.conf:/alertmanager/alertmanager.conf:ro \
            --publish 0.0.0.0:{{ prometheus_alertmanager_port }}:9093/tcp \
            {{ prometheus_alertmanager_image }} -config.file=/alertmanager/alertmanager.conf
        # iptables
        /usr/bin/docker exec --user root {{ prometheus_alertmanager_container }} /etc/docker-iptables/rules
    fi
end script

post-stop script
    if docker inspect {{ prometheus_alertmanager_container }} > /dev/null 2>&1; then
        response=$(docker inspect --format="{{ '{{' }} .State.Running {{ '}}' }}" {{ prometheus_alertmanager_container }})
        if "$response" == "true"; then
            /usr/bin/docker stop {{ prometheus_alertmanager_container }}
        fi
        /usr/bin/docker rm -fv {{ prometheus_alertmanager_container }}
    fi
end script
