# roles/base/prometheus/check/apt/templates/crontab
#!/bin/sh

/usr/local/bin/prometheus-check-apt --promfile {{ prometheus_check_apt_textfile }} --regex '{{ prometheus_check_apt_regex }}'
