# roles/base/prometheus/postgres-exporter/templates/upstart-postgres-exporter.conf
# {{ ansible_managed }}

description "Postgres Exporter"

start on filesystem and started docker

pre-start script
    if ! docker inspect {{ postgres_exporter_container }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ postgres_exporter_container }} \
            --net=host \
            --env=DATA_SOURCE_NAME="postgresql://{{ postgres_exporter_user }}:{{ postgres_exporter_pass }}@127.0.0.1:{{ postgres_exporter_port }}/postgres?sslmode=disable" \
            --volume /etc/container/{{ postgres_exporter_container }}/queries.yaml:/queries.yaml:ro \
            {{ postgres_exporter_image }} \
            --extend.query-path /queries.yaml \
            --constantLabels id={{ inventory_hostname }},group={{ postgres_exporter_group }}
    fi
end script

post-stop script
    if docker inspect {{ postgres_exporter_container }} > /dev/null 2>&1; then
        response=$(docker inspect --format="{{ '{{' }} .State.Running {{ '}}' }}" {{ postgres_exporter_container }})
        if "$response" == "true"; then
            /usr/bin/docker stop {{ postgres_exporter_container }}
        fi
        /usr/bin/docker rm -fv {{ postgres_exporter_container }}
    fi
end script