# roles/base/prometheus/iptables/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus iptables-exporter"

start on local-filesystems
stop on runlevel [!2345]
respawn
kill timeout 5

script
    exec iptables-exporter --port {{ prometheus_iptables_port }}
end script
