# roles/base/sensu/rabbitmq/templates/sensu-rabbitmq.conf
# {{ ansible_managed }}
{% set CONTAINER_NAME = sensu_prefix + '_sensu_rabbitmq' %}

description "{{ CONTAINER_NAME }} container"
author "Massimiliano Ravelli <m.ravelli@mastervoice.it>"

start on filesystem and started docker
stop on runlevel [!2345]
respawn

pre-start script
    # Container
    if ! docker inspect {{ CONTAINER_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker create \
            --name {{ CONTAINER_NAME }} \
            --hostname {{ CONTAINER_NAME }} \
            --volume {{ sensu_rabbitmq_ssl_dir }}:/etc/rabbitmq/ssl:ro \
            --publish {{ sensu_rabbitmq_port }}:5671 \
            --publish {{ sensu_rabbitmq_web_port }}:15672 \
            --env RABBITMQ_PASSWORD={{ sensu_rabbitmq_password }} \
            mtek/sensu_rabbitmq
    fi
end script

script
    /usr/bin/docker start -a {{ CONTAINER_NAME }}
end script
