# roles/base/sensu/rabbitmq/templates/sensu-rabbitmq.conf
# {{ ansible_managed }}
{% set CONTAINER_NAME = sensu_prefix + '_sensu_rabbitmq' %}

description "{{ CONTAINER_NAME }} container"
author "Massimiliano Ravelli <m.ravelli@mastervoice.it>"

start on started docker
stop on stopping docker
respawn

pre-start script
    if docker inspect {{ CONTAINER_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker rm -fv {{ CONTAINER_NAME }}
    fi
end script

script
    /usr/bin/docker run \
        --name {{ CONTAINER_NAME }} \
        --hostname {{ CONTAINER_NAME }} \
        --volume {{ sensu_rabbitmq_ssl_base_dir }}/{{ sensu_prefix }}:/etc/rabbitmq/ssl:ro \
        --publish {{ sensu_rabbitmq_port }}:5671 \
        --publish {{ sensu_rabbitmq_web_port }}:15672 \
        --env=RABBITMQ_PASSWORD={{ sensu_rabbitmq_password }} \
        mtek/sensu_rabbitmq
end script

post-stop script
    /usr/bin/docker rm -fv {{ CONTAINER_NAME }}
end script
