# roles/base/sensu/web/templates/sensu-web.conf
# Upstart: /etc/init/sensu-api.conf

description "sensu api"

start on (local-filesystems and net-device-up)
stop on [!12345]

respawn

{% if not sensu_web_run_as_root %}
setuid uchiwa
setgid uchiwa
{% endif %}
chdir /opt/uchiwa/src

script
  exec /opt/uchiwa/bin/uchiwa -c /etc/sensu/uchiwa.json -p /opt/uchiwa/src/public >> /var/log/uchiwa/info.log 2> /var/log/uchiwa/error.log
end script
