# roles/base/sensu/web/templates/sensu-web.conf
# {{ ansible_managed }}
{% set API_CONTAINER_NAME = sensu_prefix + '_sensu_api' %}
{% set CONTAINER_NAME = sensu_prefix + '_sensu_web' %}

description "{{ CONTAINER_NAME }} container"
author "Massimiliano Ravelli <m.ravelli@mastervoice.it>"

start on started {{ API_CONTAINER_NAME }}
stop on runlevel [!2345]
respawn

pre-start script
    # Container
    if ! docker inspect {{ CONTAINER_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker create \
            --name {{ CONTAINER_NAME }} \
            --hostname {{ CONTAINER_NAME }} \
            --publish {{ sensu_web_port }}:3000 \
            --link={{ API_CONTAINER_NAME }}:api \
            mtek/sensu_web
    fi
end script

script
    /usr/bin/docker start -a {{ CONTAINER_NAME }}
end script

pre-stop script
    /usr/bin/docker stop {{ CONTAINER_NAME }}
end script
