# roles/base/iptables/templates/upstart-iptables-docker.conf
# {{ ansible_managed }}

description "PRE_DOCKER iptables chain"

start on started network-interface INTERFACE=docker0

pre-start script
    sleep 1
    if iptables --check FORWARD ! -i docker0 -o docker0 -j PRE_DOCKER > /dev/null 2>&1; then
        iptables --delete FORWARD ! -i docker0 -o docker0 -j PRE_DOCKER
    fi
    iptables --insert FORWARD ! -i docker0 -o docker0 -j PRE_DOCKER
end script
