# roles/base/prometheus/check/docker/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus check-docker exporter"

start on local-filesystems
stop on runlevel [!2345]
respawn
kill timeout 5

script
    exec check-docker --port {{ prometheus_check_docker_port }}
end script
