# roles/base/prometheus/postgres/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus postgres-exporter"

start on local-filesystems
stop on runlevel [!2345]
respawn
kill timeout 5

script
    export DATA_SOURCE_NAME="{{ prometheus_postgres_data_source }}"
    exec postgres_exporter -web.listen-address=":{{ prometheus_postgres_port }}"
end script
