# roles/base/prometheus/grafana/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus grafana"

start on filesystem and started docker

pre-start script
    # prometheus_grafana_container
    if ! docker inspect {{ prometheus_grafana_container }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ prometheus_grafana_container }} \
            --volume /var/lib/container/{{ prometheus_grafana_container }}:/var/lib/grafana:rw \
            --publish 0.0.0.0:{{ prometheus_grafana_port }}:3000/tcp \
            --env DATABASE_URL=sqlite3:/data/db.sqlite3 \
            {{ prometheus_grafana_image }}
    fi
end script

post-stop script
    if docker inspect {{ prometheus_grafana_container }} > /dev/null 2>&1; then
        response=$(docker inspect --format="{{ '{{' }} .State.Running {{ '}}' }}" {{ prometheus_grafana_container }})
        if "$response" == "true"; then
            /usr/bin/docker stop {{ prometheus_grafana_container }}
        fi
        /usr/bin/docker rm -fv {{ prometheus_grafana_container }}
    fi
end script
