# roles/base/mosquitto/templates/upstart.conf
# {{ ansible_managed }}

description "Mosquitto MQTT Broker {{ mosquitto_instance  }}"
author "Roger Light <roger@atchoo.org"

start on net-device-up

respawn

exec /usr/sbin/mosquitto -c /etc/mosquitto/{{ mosquitto_instance  }}/mosquitto.conf
