# roles/base/sensu/redis/templates/sensu-redis.conf
# {{ ansible_managed }}
{% set VOLUME_NAME = sensu_prefix + '_sensu_data' %}
{% set CONTAINER_NAME = sensu_prefix + '_sensu_redis' %}

description "{{ CONTAINER_NAME }} container"
author "Massimiliano Ravelli <m.ravelli@mastervoice.it>"

start on started docker
stop on stopping docker
respawn

pre-start script
    # Remove container
    if docker inspect {{ CONTAINER_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker rm -fv {{ CONTAINER_NAME }}
    fi
    # Data volume
    if ! docker inspect {{ VOLUME_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker create \
            --name {{ VOLUME_NAME }} \
            --hostname {{ VOLUME_NAME }} \
            --volume /data \
            busybox
    fi
end script

script
    /usr/bin/docker run \
        --name {{ CONTAINER_NAME }} \
        --hostname {{ CONTAINER_NAME }} \
        --volumes-from={{ VOLUME_NAME }} \
        redis:2.8 redis-server --appendonly yes
end script

post-stop script
    /usr/bin/docker rm -fv {{ CONTAINER_NAME }}
end script
