# roles/base/sensu/redis/templates/sensu-redis.conf
# {{ ansible_managed }}
{% set VOLUME_NAME = sensu_prefix + '_sensu_data' %}
{% set CONTAINER_NAME = sensu_prefix + '_sensu_redis' %}

description "{{ CONTAINER_NAME }} container"
author "Massimiliano Ravelli <m.ravelli@mastervoice.it>"

start on filesystem and started docker
stop on runlevel [!2345]
respawn

pre-start script
    # Data volume
    if ! docker inspect {{ VOLUME_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker create \
            --name {{ VOLUME_NAME }} \
            --hostname {{ VOLUME_NAME }} \
            --volume /data \
            busybox
    fi
    # Container
    if ! docker inspect {{ CONTAINER_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker create \
            --name {{ CONTAINER_NAME }} \
            --hostname {{ CONTAINER_NAME }} \
            --volumes-from={{ VOLUME_NAME }} \
            --publish {{ sensu_redis_port }}:6379 \
            redis:2.8 redis-server --appendonly yes
    fi
end script

script
    /usr/bin/docker start -a {{ CONTAINER_NAME }}
end script

pre-stop script
    /usr/bin/docker stop {{ CONTAINER_NAME }}
end script
