# roles/base/zabbix/db/templates/zabbix-server-db.conf
# {{ ansible_managed }}
{% set CONTAINER_NAME = zabbix_prefix + '_zabbix_server_db' %}
{% set VOLUME_CONTAINER_NAME = CONTAINER_NAME + '_volume' %}

description "{{ CONTAINER_NAME }} container"
author "Massimiliano Ravelli <m.ravelli@mastervoice.it>"

start on filesystem and started {{ zabbix_prefix }}-zabbix-server-db-volume
stop on runlevel [!2345]
respawn

pre-start script
    if ! docker inspect {{ CONTAINER_NAME }} > /dev/null 2>&1; then
        docker run -d --name {{ CONTAINER_NAME }} --volumes-from {{ VOLUME_CONTAINER_NAME }} -e wal_keep_segments=20 -e allowed_networks={{ ','.join(allowed_networks) }} mtek/zabbix_db:{{ zabbix_version }}
    fi
end script

script
    /usr/bin/docker start -a {{ CONTAINER_NAME }}
end script
