# roles/base/distribution/templates/upstart-frontend.conf
# {{ ansible_managed }}

description "Docker distribution frontend service"

start on started docker and started {{ distribution_service_backend }}

pre-start script
    # distribution_container_frontend
    if ! docker inspect {{ distribution_container_frontend }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ distribution_container_frontend }} \
            --volume /etc/container/{{ distribution_container_frontend }}:/etc/nginx/conf.d:ro \
            --link {{ distribution_container_backend }}:backend \
            --publish 0.0.0.0:{{ distribution_frontend_port }}:80/tcp \
            nginx:{{ distribution_nginx_tag }}
    fi
end script

post-stop script
    if docker inspect {{ distribution_container_frontend }} > /dev/null 2>&1; then
        response=$(docker inspect --format="{{ '{{' }} .State.Running {{ '}}' }}" {{ distribution_container_frontend }})
        if "$response" == "true"; then
            /usr/bin/docker stop {{ distribution_container_frontend }}
        fi
        /usr/bin/docker rm -fv {{ distribution_container_frontend }}
    fi
end script
