# roles/base/distribution/templates/upstart-backend.conf
# {{ ansible_managed }}

description "Docker distribution backend service"

start on filesystem and started docker

pre-start script
    # distribution_container_backend
    if ! docker inspect {{ distribution_container_backend }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ distribution_container_backend }} \
            --volume /var/lib/container/{{ distribution_container_backend }}:/var/lib/registry:rw \
            --env REGISTRY_STORAGE_FILESYSTEM_ROOTDIRECTORY=/var/lib/registry \
            --publish {{ distribution_port }}:5000 \
            distribution/registry:{{ distribution_tag }}
    fi
end script

post-stop script
    if docker inspect {{ distribution_container_backend }} > /dev/null 2>&1; then
        response=$(docker inspect --format="{{ '{{' }} .State.Running {{ '}}' }}" {{ distribution_container_backend }})
        if "$response" == "true"; then
            /usr/bin/docker stop {{ distribution_container_backend }}
        fi
        /usr/bin/docker rm -fv {{ distribution_container_backend }}
    fi
end script
