# roles/base/prometheus/node/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus server"

start on local-filesystems
stop on runlevel [!2345]
respawn
kill timeout 5

script
    exec /usr/local/sbin/node_exporter \
        -auth.pass="{{ prometheus_node_pass }}" \
        -auth.user="{{ prometheus_node_user }}" \
        -collector.diskstats.ignored-devices="{{ prometheus_node_diskstats_ignored_devices }}" \
        -collector.filesystem.ignored-mount-points="{{ prometheus_node_filesystem_ignored_mount_points }}" \
        -collector.megacli.command="{{ prometheus_node_megacli_command }}}" \
        -collector.netdev.ignored-devices="{{ prometheus_node_netdev_ignored_devices }}" \
        -collector.ntp.server="{{ prometheus_node_ntp_server }}" \
        -collector.textfile.directory="{{ prometheus_node_textfile_directory }}" \
        -collectors.enabled="{{ prometheus_node_collectors_enabled}}" \
        -log.level={{ prometheus_node_log_level }} \
        -web.listen-address=":{{ prometheus_node_port }}"
end script
