# roles/base/certbot/templates/cron.daily
#!/bin/sh

certbot renew --renew-hook "/etc/cron.daily/certbot_renew_hook.sh" 2>&1 | /usr/bin/logger -t certbot
