# roles/base/docker-log/templates/upstart-elastic.conf
# {{ ansible_managed }}

description "docker-log elastic containers"

start on filesystem and started docker

pre-start script
    # docker_log_container_elastic_ambassador
    if ! docker inspect {{ docker_log_container_elastic_ambassador }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ docker_log_container_elastic_ambassador }} \
            --hostname {{ docker_log_container_elastic_ambassador }} \
            --volume /var/run/docker.sock:/var/run/docker.sock \
            cpuguy83/docker-grand-ambassador --name {{ docker_log_container_elastic }}
    fi
    # docker_log_container_elastic_data
    if ! docker inspect {{ docker_log_container_elastic_data }} > /dev/null 2>&1; then
        /usr/bin/docker create \
            --name {{ docker_log_container_elastic_data }} \
            --hostname {{ docker_log_container_elastic_data }} \
            --volume /usr/share/elasticsearch/data \
            busybox:latest
    fi
    # docker_log_container_elastic
    if ! docker inspect {{ docker_log_container_elastic }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ docker_log_container_elastic }} \
            --volumes-from {{ docker_log_container_elastic_data }} \
            elasticsearch:{{ docker_log_elasticsearch_tag }}
    fi
    # docker_log_container_elastic_cron
    if ! docker inspect {{ docker_log_container_elastic_cron }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ docker_log_container_elastic_cron }} \
            --link {{ docker_log_container_elastic }}:elasticsearch \
            madron/elasticsearch-cron:latest
    fi
end script

post-stop script
    if docker inspect {{ docker_log_container_elastic_cron }} > /dev/null 2>&1; then
        /usr/bin/docker rm -fv {{ docker_log_container_elastic_cron }}
    fi
    if docker inspect {{ docker_log_container_elastic }} > /dev/null 2>&1; then
        /usr/bin/docker rm -fv {{ docker_log_container_elastic }}
    fi
end script
