# roles/base/docker-log/templates/upstart-collector.conf
# {{ ansible_managed }}

description "docker-log collector containers"

start on filesystem and started {{ docker_log_container_elastic }}

pre-start script
    # docker_log_container_collector
    if ! docker inspect {{ docker_log_container_collector }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ docker_log_container_collector }} \
            --link {{ docker_log_container_elastic_ambassador }}:es1 \
            --volume /var/run/docker.sock:/var/run/docker.sock \
            --volume /var/lib/docker/containers:/var/lib/docker/containers \
            bprodoehl/log-collector:latest
    fi
end script

post-stop script
    if docker inspect {{ docker_log_container_collector }} > /dev/null 2>&1; then
        /usr/bin/docker rm -fv {{ docker_log_container_collector }}
    fi
end script
