# roles/base/docker-log/templates/upstart-kibana.conf
# {{ ansible_managed }}

description "docker-log kibana containers"

start on filesystem and started {{ docker_log_container_elastic }}

pre-start script
    # docker_log_container_kibana
    if ! docker inspect {{ docker_log_container_kibana }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ docker_log_container_kibana }} \
            --link {{ docker_log_container_elastic_ambassador }}:elasticsearch \
            --publish {{ docker_log_kibana_port }}:5601 \
            madron/kibana:latest
    fi
end script

post-stop script
    if docker inspect {{ docker_log_container_kibana }} > /dev/null 2>&1; then
        /usr/bin/docker rm -fv {{ docker_log_container_kibana }}
    fi
end script
