# roles/base/prometheus/server/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus server"

start on filesystem and started docker

pre-start script
    # prometheus_server_container
    if ! docker inspect {{ prometheus_server_container }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ prometheus_server_container }} \
            --volume /etc/container/{{ prometheus_server_container }}/prometheus.yml:/etc/prometheus/prometheus.yml:ro \
            --volume /etc/container/{{ prometheus_server_container }}/file_sd_configs:/etc/prometheus/file_sd_configs:ro \
            --volume /etc/container/{{ prometheus_server_container }}/rule_files:/etc/prometheus/rule_files:ro \
            --volume /var/lib/container/{{ prometheus_server_container }}:/prometheus:rw \
            --publish 0.0.0.0:{{ prometheus_server_port }}:9090/tcp \
            --env=SAMPLEVAR={{ prometheus_server_container }} \
            {{ prometheus_server_image }} \
{% if prometheus_server_alertmanager_url %}
            -alertmanager.url {{ prometheus_server_alertmanager_url }} \
{% endif %}
            -config.file=/etc/prometheus/prometheus.yml \
            -storage.local.path=/prometheus \
            -storage.local.retention={{ prometheus_server_storage_local_retention }} \
            -web.console.libraries=/etc/prometheus/console_libraries \
            -web.console.templates=/etc/prometheus/consoles
    fi
end script

post-stop script
    if docker inspect {{ prometheus_server_container }} > /dev/null 2>&1; then
        response=$(docker inspect --format="{{ '{{' }} .State.Running {{ '}}' }}" {{ prometheus_server_container }})
        if "$response" == "true"; then
            /usr/bin/docker stop {{ prometheus_server_container }}
        fi
        /usr/bin/docker rm -fv {{ prometheus_server_container }}
    fi
end script
