# roles/base/consul/agent/templates/upstart.conf
# {{ ansible_managed }}

description "Consul agent"

start on local-filesystems
stop on runlevel [!2345]
respawn
kill timeout 5
setuid nobody
setgid nogroup

script
    exec /usr/local/sbin/consul agent -config-file /etc/consul/consul.json
end script
