# roles/base/cadvisor/templates/systemd.service
[Unit]
Description=Cadvisor
After=docker.service

[Service]
Restart=always
RestartSec=10
ExecStartPre=/usr/local/bin/docker-compose -p {{ project }} -f /etc/init/{{ service }}.yml down -v
ExecStart=/usr/local/bin/docker-compose -p {{ project }} -f /etc/init/{{ service }}.yml up
ExecStop=/usr/local/bin/docker-compose -p {{ project }} -f /etc/init/{{ service }}.yml down -v

[Install]
WantedBy=default.target
