# roles/base/prometheus/pushgateway/templates/upstart.conf
# {{ ansible_managed }}

description "Prometheus pushgateway"

start on filesystem and started docker

pre-start script
    # prometheus_pushgateway_container
    if ! docker inspect {{ prometheus_pushgateway_container }} > /dev/null 2>&1; then
        /usr/bin/docker run \
            -d \
            --restart="always" \
            --name {{ prometheus_pushgateway_container }} \
            --volume /var/lib/container/{{ prometheus_pushgateway_container }}:/data:rw \
            --publish 0.0.0.0:{{ prometheus_pushgateway_port }}:9091/tcp \
            {{ prometheus_pushgateway_image }} \
            -persistence.file /data/persistence \
            -log.level {{ prometheus_pushgateway_log_level }}
    fi
end script

post-stop script
    if docker inspect {{ prometheus_pushgateway_container }} > /dev/null 2>&1; then
        response=$(docker inspect --format="{{ '{{' }} .State.Running {{ '}}' }}" {{ prometheus_pushgateway_container }})
        if "$response" == "true"; then
            /usr/bin/docker stop {{ prometheus_pushgateway_container }}
        fi
        /usr/bin/docker rm -fv {{ prometheus_pushgateway_container }}
    fi
end script
