# roles/base/sensu/server/templates/sensu-server.conf
# Upstart: /etc/init/sensu-server.conf

description "sensu server"

start on (local-filesystems and net-device-up)
stop on [!12345]

respawn

setuid sensu
setgid sensu

script
  CONFIG_FILE=/etc/sensu/server.json
  CONFIG_DIR=/etc/sensu/server-conf.d
  EXTENSION_DIR=/etc/sensu/enabled-extensions
  LOGFILE=/var/log/sensu/sensu-server.log

  exec /opt/sensu/bin/sensu-server -c $CONFIG_FILE -d $CONFIG_DIR -e $EXTENSION_DIR -l $LOGFILE
end script
