# roles/base/sensu/server/templates/sensu-server.conf
# {{ ansible_managed }}
{% set RABBITMQ_CONTAINER_NAME = sensu_prefix + '_sensu_rabbitmq' %}
{% set REDIS_CONTAINER_NAME = sensu_prefix + '_sensu_redis' %}
{% set CONTAINER_NAME = sensu_prefix + '_sensu_server' %}

description "{{ CONTAINER_NAME }} container"
author "Massimiliano Ravelli <m.ravelli@mastervoice.it>"

start on started {{ RABBITMQ_CONTAINER_NAME }} and started {{ REDIS_CONTAINER_NAME }}
stop on runlevel [!2345]
respawn

pre-start script
    if docker inspect {{ CONTAINER_NAME }} > /dev/null 2>&1; then
        /usr/bin/docker rm -fv {{ CONTAINER_NAME }}
    fi
end script

script
    /usr/bin/docker run \
        --name {{ CONTAINER_NAME }} \
        --hostname {{ CONTAINER_NAME }} \
        --volume {{ sensu_server_etc_base_dir }}/{{ sensu_prefix }}/ssl:/etc/sensu/ssl:ro \
        --volume {{ sensu_server_etc_base_dir }}/{{ sensu_prefix }}/conf.d:/etc/sensu/conf.d:ro \
        --link={{ RABBITMQ_CONTAINER_NAME }}:rabbitmq \
        --link={{ REDIS_CONTAINER_NAME }}:redis \
        --env=RABBITMQ_PASSWORD={{ sensu_rabbitmq_password }} \
        mtek/sensu_server
end script

pre-stop script
    /usr/bin/docker rm -fv {{ CONTAINER_NAME }}
end script
