# roles/base/sensu/server/templates/sensu-api.conf
# Upstart: /etc/init/sensu-api.conf

description "sensu api"

start on (local-filesystems and net-device-up)
stop on [!12345]

respawn

setuid sensu
setgid sensu

script
  CONFIG_FILE=/etc/sensu/api.json
  LOGFILE=/var/log/sensu/sensu-api.log

  exec /opt/sensu/bin/sensu-api -c $CONFIG_FILE -l $LOGFILE
end script
